# deploy/com.b9.xcode-auto-approve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.b9.xcode-auto-approve</string>

    <key>ProgramArguments</key>
    <array>
        <string>__BINARY_PATH__</string>
    </array>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/xcode-auto-approve.log</string>

    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/xcode-auto-approve.log</string>
</dict>
</plist>
